# hermes-agent/lihi-coach/systemd/hermes-lihi.service
[Unit]
# Hermes instance for Lihi's CMO coach ("Noga") — dedicated user, narrower
# sandbox than the main hermes unit: ONLY the marketing repo is bound in
# (no /home/dbexpertai/code, no Claude credentials, no claude CLI — the
# coach converses and files documents; it does not orchestrate dev work).
Description=Hermes gateway — Lihi's CMO coach (hardened, user=hermes-lihi)
Documentation=file:///home/dbexpertai/code/jetsonlocalai/hermes-agent/lihi-coach/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hermes-lihi
Group=hermes-lihi
WorkingDirectory=/home/hermes-lihi/.hermes

Environment="PATH=/home/hermes-lihi/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/home/hermes-lihi"
Environment="HERMES_HOME=/home/hermes-lihi/.hermes"

ExecStart=/home/hermes-lihi/.local/bin/hermes gateway run

Restart=on-failure
RestartSec=10
TimeoutStopSec=240

# --- Filesystem isolation (see main hermes.service for rationale) ---------
# Bind in ONLY the coach's home and the marketing workspace. David's other
# repos, credentials, keys: absent from this namespace entirely.
ProtectHome=tmpfs
BindPaths=/home/hermes-lihi
BindPaths=/home/dbexpertai/code/marketing-liram-heshev

ProtectSystem=strict
ReadWritePaths=/home/hermes-lihi /home/dbexpertai/code/marketing-liram-heshev
PrivateTmp=true

# --- Kernel surface hardening ---------------------------------------------
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true

# --- Privilege boundaries --------------------------------------------------
NoNewPrivileges=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
CapabilityBoundingSet=
AmbientCapabilities=

# Localhost services the coach uses: jetson-stt (127.0.0.1:11436).
# Outbound HTTPS needed for Telegram polling + LLM inference.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

OOMScoreAdjust=-50

[Install]
WantedBy=multi-user.target
